# deploy/feishu-blog-bot.service
[Unit]
Description=Feishu Blog Agent Bot (WebSocket)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/AI-Sticker-Ecommerce
ExecStart=/root/miniconda3/envs/sticker/bin/python feishu_blog_app.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

EnvironmentFile=/root/AI-Sticker-Ecommerce/.env

StandardOutput=append:/var/log/feishu-blog-bot/stdout.log
StandardError=append:/var/log/feishu-blog-bot/stderr.log

[Install]
WantedBy=multi-user.target
