# deploy/feishu-sticker-bot.service
[Unit]
Description=Feishu Sticker Pack Bot (WebSocket)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/root/AI-Sticker-Ecommerce
ExecStart=/root/miniconda3/envs/sticker/bin/python feishu_chat_app.py
Restart=always
RestartSec=5

EnvironmentFile=/root/AI-Sticker-Ecommerce/.env

StandardOutput=append:/var/log/feishu-sticker-bot/stdout.log
StandardError=append:/var/log/feishu-sticker-bot/stderr.log

[Install]
WantedBy=multi-user.target
