# deploy/feishu-blog-bot.service
[Unit]
Description=Feishu Blog Agent Bot (WebSocket)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/ai-sticker-ecommerce
ExecStart=/opt/ai-sticker-ecommerce/venv/bin/python feishu_blog_app.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

EnvironmentFile=/opt/ai-sticker-ecommerce/.env

StandardOutput=append:/var/log/feishu-blog-bot/stdout.log
StandardError=append:/var/log/feishu-blog-bot/stderr.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/ai-sticker-ecommerce/data /opt/ai-sticker-ecommerce/output /opt/ai-sticker-ecommerce/logs /var/log/feishu-blog-bot

[Install]
WantedBy=multi-user.target
